# gpu_stats.timer
[Unit]
Description=GPU stats every 5 minutes

[Timer]
OnBootSec=5m
OnUnitInactiveSec=5m

[Install]
WantedBy=basic.target

