# gpu_stats.service
[Unit]
Description=GPU stats service

[Service]
Type=oneshot
StandardOutput=null
ExecStart=/usr/bin/gpu_stats %h
